# rootfs/etc/s6-overlay/s6-rc.d/neservice/run
#!/command/execlineb -P
# shellcheck shell=sh
fdmove -c 2 1
/usr/local/netextender/NEService
